# crontab
* * * * * env OTP_SECRET=SYZRVQURZIXIUOIR python3 /app/src/main.py >> /var/log/cron.log 2>&1


